# crontab
@reboot screen -dm -S emacs-daemon ssh-agent /home/bvk/bin/emacs --fg-daemon
